# systemd/geometryos-directive-daemon.service
[Unit]
Description=Geometry OS Directive Polling Daemon
Documentation=https://github.com/geometry-os/geometry-os
After=network.target

[Service]
Type=simple
User=%USER%
WorkingDirectory=%WORKDIR%
ExecStart=%PYTHON% %WORKDIR%/scripts/directive_daemon.py
Restart=on-failure
RestartSec=5s

# Environment configuration
Environment="WP_URL=http://localhost:8080/ai-publisher.php"
Environment="POLL_INTERVAL=5"
Environment="USE_LLM=false"
Environment="LLM_ENDPOINT=http://localhost:11434/api/generate"
Environment="LLM_MODEL=llama3.2:latest"

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%WORKDIR%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=geometryos-directive-daemon

[Install]
WantedBy=default.target
